# ops/systemd/systemtrader-runner.service
[Unit]
Description=SystemTrader Hybrid scheduled scanner runner
After=systemtrader-browser.service systemtrader-telegram-relay.service
Wants=systemtrader-browser.service

[Service]
Type=oneshot
User=systemtrader
Group=systemtrader
EnvironmentFile=/etc/systemtrader/systemtrader.env
WorkingDirectory=/opt/systemtrader
ExecStart=/usr/bin/node /opt/systemtrader/ops/scanner-runner.js
